# deps/nsproxy.service
[Unit]
Description=nsproxy

[Service]
ExecStart=/usr/bin/rkt run \
	--net=host \
        --volume redis,kind=host,source=/tmp/dump.rdb \
        unixvoid.com/nsproxy

[Install]
WantedBy=multi-user.target
